# bot.service
[Unit]
Description=XAU/USD HFT Bot - Connection Core
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory={{PROJECT_DIR}}
ExecStart={{PYTHON_PATH}} main.py
Environment=PYTHONPATH={{PROJECT_DIR}}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# High-Performance Tuning
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=99
Nice=-10

[Install]
WantedBy=multi-user.target
